# systemd/batcaplog.service
[Unit]
Description=Log battery capacity

[Service]
Type=oneshot
ExecStart=/usr/bin/batcaplog
